# cron/root.crontab
@hourly /home/el/.local/scripts/downloadupdate
